# files/supervisor/wso2-greg.ini
[program:wso2-greg]
user = root
autostart = true
command = /opt/wso2/greg/bin/wso2server.sh
stdout_logfile = /var/log/wso2/greg/common.log
stderr_logfile = /var/log/wso2/greg/error.log

